# profile.d/dotnet.sh
export PATH="${HOME}/.dotnet:${PATH}:${HOME}/.dotnet/tools"
export DOTNET_ROOT="${HOME}/.dotnet"

